# metricool-mcp.service
# /etc/systemd/system/metricool-mcp.service
#
# Systemd unit for the Metricool MCP server.
#
# Install:
#   sudo cp metricool-mcp.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now metricool-mcp
#
# Logs:
#   sudo journalctl -u metricool-mcp -f

[Unit]
Description=Metricool MCP Server
After=network.target

[Service]
Type=simple
User=metricool
Group=metricool
WorkingDirectory=/opt/metricool-mcp

# --- Environment ---
# Single-user mode: set the API key here.
# Multi-user mode: leave blank — clients send their key via Authorization: Bearer.
EnvironmentFile=-/opt/metricool-mcp/.env
Environment=METRICOOL_BASE_URL=https://app.metricool.com
Environment=HOST=127.0.0.1
Environment=PORT=8000

# --- Process ---
ExecStart=/opt/metricool-mcp/venv/bin/uvicorn server:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=on-failure
RestartSec=5s

# --- Security hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/metricool-mcp

[Install]
WantedBy=multi-user.target
